# meta-openvision/recipes-connectivity/bluez5/files/bluez-alsa.service
[Unit]
Description=Start the bluez-alsa service
After=bluetooth.service
After=systemd-remount-fs.service

[Service]
ExecStart=${bindir}/bluealsa

[Install]
WantedBy=bluetooth.target
